# crontab
*/30 * * * * /usr/bin/unison ~/.unison/ratchet-to-springer.prf &> unison.log
